# service/hdmi-display.service
[Unit]\nDescription=HDMI Display Service\n[Service]\nExecStart=/usr/bin/somecommand\n[Install]\nWantedBy=multi-user.target\n