# service/hdmi-display.service
[Unit]
Description=HDMI Simple Display Service
After=display-manager.service
Wants=display-manager.service

[Service]
Type=simple
# Set TEST_PATTERN_PATH to customize the test pattern image
Environment="TEST_PATTERN_PATH=/usr/share/hdmi-display/test-pattern.png"
ExecStart=/bin/sh -c 'if [ -f "$TEST_PATTERN_PATH" ]; then exec /usr/bin/hdmi_simple_display --test-pattern="$TEST_PATTERN_PATH"; else exec /usr/bin/hdmi_simple_display; fi'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target